# .claude/skills/claude/settings/crontab
0 0 * * * bash /Users/duck/Documents/Work/DuckyWorld/.claude/skills/claude/auto-update.sh
