# .claude/skills/claude/settings/crontab
# Claude Code 자동 업데이트 - 매일 00시
0 0 * * * bash __SKILL_DIR__/auto-update.sh
